# deploy/ansible/templates/hiretrack-sync-full-refresh.timer.j2
[Unit]
Description=Run HireTrack weekly full refresh

[Timer]
OnCalendar=Sun 03:00:00
Persistent=true
RandomizedDelaySec={{ hiretrack_full_refresh_randomized_delay_sec }}
Unit=hiretrack-sync-full-refresh.service

[Install]
WantedBy=timers.target
