# deploy/systemd/hiretrack-sync-full-refresh.timer
[Unit]
Description=Run HireTrack weekly full refresh

[Timer]
OnCalendar=Sun 03:00:00
Persistent=true
RandomizedDelaySec=900
Unit=hiretrack-sync-full-refresh.service

[Install]
WantedBy=timers.target
